# fan-daemon.service
[Unit]
Description=Fan control daemon

[Service]
Type=simple
ExecStart=/usr/local/bin/fan-daemon --interval 2 --fans 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target